# 9.1/nailgun.ini
[program:nailgun]
command=/usr/bin/manage.py run
