# 9.1/nailgun.ini
[program:nailgun]
command=/usr/bin/manage.py run --fake-tasks --fake-tasks-tick-count=80 --fake-tasks-tick-interval=1 --authentication-method="none"
stdout_logfile=/var/log/nailgun/app.log
redirect_stderr=true
