# 9.1/postgresql.ini
[program:postgresql]
command=/usr/bin/postgres -D /var/lib/pgsql/data -p 5432
user=postgres
